# scripts/com.nickschneble.looper-custodian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.nickschneble.looper-custodian</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>/Users/nickschneble/Developer/Repos/agents-of-shield-if-shield-is-ai/scripts/looper-custodian-cron.sh</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Weekday</key>
		<integer>1</integer>
		<key>Hour</key>
		<integer>9</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>StandardOutPath</key>
	<string>/Users/nickschneble/Developer/Repos/agents-of-shield-if-shield-is-ai/local/custodian/launchd.out.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/nickschneble/Developer/Repos/agents-of-shield-if-shield-is-ai/local/custodian/launchd.err.log</string>
	<key>RunAtLoad</key>
	<false/>
	<key>ProcessType</key>
	<string>Background</string>
</dict>
</plist>
